# haipeng/roles/slave/template/spring-boot.service
[Unit]
Description=Spring-boot

[Service]
ExecStart=java -jar /root/springboot_vue-0.0.1-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target